# task/.config/task/sync/taskwarrior_git_push.timer
[Unit]
Description=Run Taskwarrior Git Push Script Periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
